# systemd/weatheraqi.target
[Unit]
Description=Weather AQI Target
StopWhenUnneeded=yes
